# sofan-ai-fresh/deploy/systemd/hermes-gateway.service
[Unit]
Description=Hermes Agent Gateway - Messaging Platform Integration
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=%h/.local/bin/uv run --directory %h/.hermes/hermes-agent python -m hermes_cli.main gateway run
WorkingDirectory=%h/.hermes
Environment=HERMES_HOME=%h/.hermes
Restart=always
RestartSec=5
RestartForceExitStatus=75
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=90
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
